# .config/systemd/user/miner.service
[Unit]
Description=Start mining
After=network.target

[Service]
Type=idle
ExecStart=/home/farmer/miners/start.sh

[Install]
WantedBy=default.target
